# install/strava-dashboard.service
[Unit]
Description=Strava E-Paper Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/strava-dashboard
Restart=always
RestartSec=30
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
